# scada_plc/docker/dockerfiles/manual/shared/supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 1MB
logfile_backups = 0

[program:browse]
command = /usr/local/bin/poc_browse.py

